# supervisord.conf
[supervisord]
nodaemon=false

[program:mariadb]
command=/usr/bin/mysqld_safe --datadir='/var/lib/mysql' --port=3306 --skip-networking=0